# projects/Legacy/PureDarwinSettings/System/Library/LaunchDaemons/org.ntp.ntpd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>org.ntp.ntpd</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/libexec/ntpd-wrapper</string>
    </array>
    <key>KeepAlive</key>
    <dict>
      <key>PathState</key>
      <dict>
	<key>/private/etc/ntp.conf</key>
	<true/>
      </dict>
    </dict>
<!--
	Un-comment the following lines to enable overloading launchd conf file schema...
-->    
<!--    
    <key>SHAuthorizationRight</key>
    <string>system.preferences</string>
-->    
  </dict>
</plist>

