# src/Legacy/PureDarwinSettings/System/Library/LaunchDaemons/com.apple.securityd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.SecurityServer</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/sbin/securityd</string>
		<string>-i</string>
		<string>-t</string>
		<string>10</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.apple.SecurityServer</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>LaunchOnlyOnce</key>
	<true/>
</dict>
</plist>

