# src/libexec/syslog/aslmanager.tproj/com.apple.aslmanager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.aslmanager</string>
	<key>Program</key>
	<string>/usr/sbin/aslmanager</string>
	<key>ProgramArguments</key>
	<array>
		<string>aslmanager</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.apple.aslmanager</key>
		<true/>
	</dict>
	<key>EnableTransactions</key>
	<true/>
	<key>EnablePressuredExit</key>
	<true/>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
</dict>
</plist>
