# projects/PureDarwinSettings/System/Library/LaunchDaemons/com.apple.hidd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.apple.hidd</string>
    <key>MachServices</key>
    <dict>
            <key>com.apple.iohideventsystem</key>
            <true/>
    </dict>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/hidd</string>
	</array>
	<key>Disabled</key>
	<true/>
</dict>
</plist>
